# deployment/etc/systemd/system/lights-watcher.path
[Path]
PathModified=/home/lights/lights

[Install]
WantedBy=multi-user.target
